# dmenu-quicklaunch-srv.service
[Unit]
Description=Dmenu Desktop Launcher Server

[Service]
ExecStart=/usr/bin/env dmenu-quicklaunch-srv

[Install]
WantedBy=default.target
